# simulator-auto-upgrade.timer
[Unit]
Description=Timer for auto updating

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
